# bgutil-ytdlp-pot-provider-git/proxy-to-bgutil-ytdlp-pot-provider.socket
[Socket]
ListenStream=4416

[Install]
WantedBy=sockets.target
